# preordain/scraper/crontab
# Place the path to the directory in the commands
@daily python -m preordain.scraper.price_scrape
@weekly python -m preordain.scraper.sale_scrape
