# scraper/crontab
# Place the path to the directory in the commands
@daily python -m /app.price_scrape
@weekly python -m /app.sale_scrape
